# nodes/ogma.podman/quadlets/openbao.container
# Podman Quadlet — pure git-managed unit, deployed only by the gitops-pull timer
# (ansible/roles/gitops_pull) into /etc/containers/systemd/. No Ansible role templates this
# file. TLS material and the mesh firewall rule are host-level and provisioned separately by
# ansible/roles/openbao — nothing sensitive lives in this repo.
[Unit]
Description=OpenBao server
After=network-online.target

[Container]
Image=ghcr.io/openbao/openbao:2
ContainerName=openbao
Volume=openbao.volume:/openbao/data
Volume=/etc/openbao/config:/openbao/config:ro,Z
Volume=/etc/openbao/tls:/openbao/tls:ro,Z
PublishPort=8200:8200
Exec=bao server -config=/openbao/config/openbao.hcl

[Service]
Restart=always

[Install]
WantedBy=multi-user.target
