# nodes/ogma.podman/quadlets/traefik.container
# Podman Quadlet — pure git-managed unit, deployed only by the gitops-pull timer
# (ansible/roles/gitops_pull) into /etc/containers/systemd/. No Ansible role templates this
# file. The Bunny API key (for DNS-01) is host-level and provisioned separately by
# ansible/roles/pocketid — nothing sensitive lives in this repo. Routes to other containers
# purely via Docker-compatible labels read off the podman socket — exposedByDefault is off, so
# only containers carrying traefik.enable=true (currently just pocketid) are ever routed to.
[Unit]
Description=Traefik (ogma public edge)
After=network-online.target

[Container]
Image=docker.io/library/traefik:v3
ContainerName=traefik
# Rootful podman's bridge networking masquerades even externally-forwarded traffic on published
# ports — every connection appeared to Traefik as coming from the bridge gateway (10.88.0.1),
# breaking the vault-tailnet ipAllowList (and any other source-IP-based logic). Host networking
# skips the bridge/NAT layer entirely, so Traefik sees real client IPs. No PublishPort needed —
# entryPoint addresses (see traefik.yml) bind directly on the host.
Network=host
Volume=/etc/traefik/config/traefik.yml:/etc/traefik/traefik.yml:ro,Z
Volume=/run/podman/podman.sock:/run/podman/podman.sock:ro
Volume=traefik-acme.volume:/letsencrypt
EnvironmentFile=/etc/traefik/env
Exec=traefik --configFile=/etc/traefik/traefik.yml
# The socket keeps its native container_runtime_t label (relabeling it with :Z/:z only ever gets
# it to container_file_t, which podman's own SELinux policy doesn't grant connectto on). Instead
# put Traefik's own process in the same domain podman.service already runs in — the documented
# fix for mounting the podman/docker socket into a container:
# https://github.com/containers/podman/discussions/9860
SecurityLabelType=container_runtime_t

[Service]
Restart=always

[Install]
WantedBy=multi-user.target
