# ansible/roles/gitops_pull/templates/futhark-gitops-pull.service.j2
[Unit]
Description=Futhark GitOps pull-reconcile ({{ node.hostname }})

[Service]
Type=oneshot
ExecStart=/usr/local/bin/futhark-gitops-pull.sh
